# deploy/gunicorn.service.example
# systemd unit for Gunicorn (Nana's Bites backend)
# Place as /etc/systemd/system/nanas-bites.service

[Unit]
Description=Nana's Bites Gunicorn
After=network.target postgresql.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/nanas_bites/backend
EnvironmentFile=/var/www/nanas_bites/backend/.env
ExecStart=/var/www/nanas_bites/venv/bin/gunicorn config.wsgi:application -c gunicorn.conf.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
